# profile.d/gopass.sh
export GOPASS_GPG_OPTS="--no-autostart"
#alias gopass="$(which gopass) $@ 2> /dev/null"
