# profile.d/mysql-client.sh
export PATH="/home/linuxbrew/.linuxbrew/opt/mysql-client/bin:$PATH"
